# sonarr.conf
[program:sonarr]
user = nobody
command = /usr/bin/sonarr -nobrowser -data=/config